# containers/dns/supervisord.conf
; supervisord program list for the DNS container (issue #248).
; Replaces the bare `CMD ["named", ...]` so the in-process Wazuh agent
; can run alongside named under a single supervisor.

[program:named]
command=/usr/sbin/named -g -c /etc/bind/named.conf -u bind
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wazuh-agent]
; In-process Wazuh agent. Bootstrap script registers with the manager
; (with `agent-auth -F 1` to take over from any prior sidecar
; registration of the same name), then execs `tail -F` on
; /var/ossec/logs/ossec.log so supervisord can manage lifecycle. Env
; vars (WAZUH_MANAGER, AGENT_NAME, LOG_PATHS, LOG_FORMAT) come from the
; container env block in docker-compose.yml.
command=/opt/aptl/wazuh/wazuh-agent.sh
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=WAZUH_AGENT_TEMPLATE="/opt/aptl/wazuh/ossec.conf.template"
