# containers/fileshare/supervisord.conf
; supervisord program list for the fileshare container (issue #248).
; Replaces the previous `exec smbd --foreground` at the end of
; setup-shares.sh; supervisord now owns samba + the in-process Wazuh
; agent + rsyslog (which the previous setup-shares.sh started in the
; background).

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:samba]
command=/usr/sbin/smbd --foreground --no-process-group --debug-stdout
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wazuh-agent]
; In-process Wazuh agent. See dns/supervisord.conf for protocol details.
command=/opt/aptl/wazuh/wazuh-agent.sh
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=WAZUH_AGENT_TEMPLATE="/opt/aptl/wazuh/ossec.conf.template"
