# containers/webapp/supervisord.conf
; supervisord program list for the webapp container (issue #248).
; Replaces the previous `exec gunicorn` at the end of entrypoint.sh;
; supervisord now owns gunicorn, rsyslog, and the in-process Wazuh
; agent.

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
command=/usr/local/bin/gunicorn --bind 0.0.0.0:8080 --workers 2 --access-logfile /var/log/gunicorn/access.log --error-logfile - app:app
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wazuh-agent]
; In-process Wazuh agent. See dns/supervisord.conf for protocol details.
command=/opt/aptl/wazuh/wazuh-agent.sh
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=WAZUH_AGENT_TEMPLATE="/opt/aptl/wazuh/ossec.conf.template"
